# systemd/rdevents-backup.service
[Unit]
Description=RDevents scheduled backup

[Service]
Type=oneshot
ExecStart=/opt/rdevents/ops/scheduled-backup.sh
