# host-systemd/rdevents-system-report-refresh.service
[Unit]
Description=RDEvents System Report Generator
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/rdevents/bin/build-system-report.sh
WorkingDirectory=/opt/rdevents
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rdevents-system-report
User=root
Group=root
TimeoutStartSec=180
