# host-systemd/rdevents-system-report-refresh.path
[Unit]
Description=RDEvents System Report Refresh Watcher
Documentation=
After=network.target

[Path]
PathExists=/opt/rdevents/runtime/control/system-report-refresh-request.json
Unit=rdevents-system-report-refresh.service

[Install]
WantedBy=multi-user.target
